# files/etc/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command=/usr/sbin/php-fpm5.6 -F
autostart=true
autorestart=true
startretries=3
